# bin/derpme-docker.service
[Unit]
Description=Service for DerpMe.
# After=multi-user.target
After=multi-user.target network.target redis-server.service

[Service]
User=klpanagi
Type=idle
ExecStart=/usr/bin/docker run --rm --name derpme-server --net host derpme-server
ExecStop=/usr/bin/docker container rm derpme-server --force
StandardOutput=inherit
StandardError=inherit
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
