# bin/derpme.service
[Unit]
Description=Service for DerpMe.
# After=multi-user.target
After=multi-user.target network.target redis-server.service

[Service]
User=pi
Type=idle
ExecStart=derp_me
StandardOutput=inherit
StandardError=inherit
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
